# scripts/project-forge-expand.service
[Unit]
Description=Project Forge Horizontal Expansion
After=network.target

[Service]
Type=oneshot
User=claude
WorkingDirectory=/opt/project-forge
Environment=FORGE_DB_PATH=/opt/project-forge/data/forge.db
EnvironmentFile=-/opt/project-forge/.env
ExecStart=/opt/project-forge/scripts/expand-hourly.sh
